# configs/auto-update.timer
[Unit]
Description=Pi OBD Meter auto-update check

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
