# configs/kiosk.service
[Unit]
Description=OBD Meter Kiosk Display (cage + Chromium)
After=pi-obd-meter.service
Wants=pi-obd-meter.service

[Service]
Type=simple
User=laurel
PAMName=login
TTYPath=/dev/tty7
Environment=WLR_LIBINPUT_NO_DEVICES=1
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/bin/bash /opt/pi-obd-meter/configs/kiosk.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
